# config/systemd/user/water.timer
[Unit]
Description=Drink some water, fool

[Timer]
OnBootSec=15m
OnUnitActiveSec=15m
Unit=water.service

[Install]
WantedBy=timers.target
